# lxc-manager.service
[Unit]
Description=LXC Simple Manager
Documentation=https://github.com/joseprupi/lxc-simple-manager
After=network.target lxc.service lxc-net.service
Wants=lxc.service lxc-net.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lxc_manager
ExecStart=/opt/lxc_manager/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8001 --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening (limited since we need root for LXC)
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
